# source/app/misc/ioc.service
[Unit]
Description=IOC

[Service]
ExecStart=/usr/bin/bash /opt/run_ioc.sh
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
